# code/automation/automation.service
[Unit]
Descriptin= install httpd with ansible.

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/automate.sh

[Install]
WantedBy=multi-user.target
